# hive/opt/repomirror/etc/cron.hourly/hive-updaterepo
#!/usr/bin/env bash

(/hive/opt/repomirror/updaterepo 2>&1) | tee /var/log/hive-repo-sync.log
